# dot_config/systemd/user/rsync_backup.timer.tmpl
{{ if .personal -}}
[Unit]
Description=Sync user data to homenas

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
{{ end -}}
